# deploy/deb/systemd/zfs_exporter.service
[Unit]
Description=Prometheus ZFS Exporter
Documentation=https://github.com/donaldgifford/zfs_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zfs_exporter
Group=zfs_exporter
EnvironmentFile=-/etc/default/zfs_exporter
ExecStart=/usr/bin/zfs_exporter $ZFS_EXPORTER_OPTS
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# PrivateDevices must be "no" -- zpool/zfs commands require /dev/zfs
# to communicate with the ZFS kernel module. PrivateDevices=yes hides
# /dev/zfs, causing zpool commands to hang indefinitely.
PrivateDevices=no
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadOnlyPaths=/
ReadWritePaths=/proc/spl/kstat/zfs
DeviceAllow=/dev/zfs rw

# Allow access to zpool/zfs binaries and systemctl
ExecPaths=/usr/bin /usr/sbin /bin /sbin

[Install]
WantedBy=multi-user.target
